# telegram-bot.service
[Unit]
Description=Claw Telegram Bot (voice-enabled)
After=network.target openclaw-router.service

[Service]
ExecStart=/usr/bin/python3 /home/ubuntu/telegram_bot.py
Restart=always
RestartSec=5
WorkingDirectory=/home/ubuntu

# ── Required ──────────────────────────────────────────────────────────────────
Environment=TELEGRAM_BOT_TOKEN=8233548348:AAF-MfrapA4msggPkzuwy75wBaMNrdBsvjQ

# ── Optional — Groq API for fast cloud Whisper (recommended) ──────────────────
# Get free key at: https://console.groq.com (no credit card needed)
# Provides whisper-large-v3 at near-realtime speed
# Without this, falls back to local faster-whisper/small (slower but offline)
#Environment=GROQ_API_KEY=gsk_xxxxxxxxxxxxxxxxxxxx

# ── Voice model (used only if GROQ_API_KEY not set) ───────────────────────────
# Options: tiny (fast, poor), base, small (recommended), medium (best, slow)
Environment=WHISPER_MODEL=small

# ── Router config ─────────────────────────────────────────────────────────────
Environment=ROUTER_URL=http://localhost:4097/v1/chat/completions
Environment=ROUTER_MODEL=opencode/minimax-m2.5-free

# ── Access control (optional) ─────────────────────────────────────────────────
# Comma-separated Telegram chat IDs that can use the bot (empty = allow all)
# Find your chat ID by messaging @userinfobot
#Environment=ALLOWED_CHAT_IDS=123456789,987654321

[Install]
WantedBy=default.target
